# fw-log-query.service
[Unit]
Description=Firewall Log Query System - DuckDB Optimized
Documentation=https://github.com/your-org/fw-log-query
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/fw_log_query

# 启动命令
ExecStart=/opt/fw_log_query/fw_log_query \
    -logdir=/data/sangfor_fw_log \
    -db=/data/index/fw_logs.duckdb \
    -port=8080 \
    -workers=4

# 重启策略
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=infinity

# 性能优化
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=2

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/sangfor_fw_log /data/index

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fw-log-query

[Install]
WantedBy=multi-user.target
